# deploy/dash.conf
description "dash"

start on (filesystem)
stop on runlevel [016]

respawn
console log
#setuid user
#setgid user
chdir /home/ian/tyler

exec /home/ian/tyler/deploy/run.sh
